# scripts/archive/openclaw_retired/deploy/lima-openclaw.service
[Unit]
Description=LiMa OpenClaw Gateway (WeChat light, LiMa brain)
After=network.target lima-router.service
Wants=lima-router.service

[Service]
Type=simple
User=root
OOMScoreAdjust=-900
WorkingDirectory=/opt/lima-router
Environment=HOME=/root
Environment=PATH=/root/.nvm/versions/node/v22.22.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/opt/lima-router/scripts/openclaw_gateway_start.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
