# deploy/jdcloud/lima-probe-push.service
[Unit]
Description=LiMa Probe Result Push to Main VPS
After=network.target network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/lima-probe
ExecStart=/usr/bin/python3 /opt/lima-probe/push_probe_results.py
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/lima-probe/.env

# Security
NoNewPrivileges=yes
PrivateTmp=yes
