# deploy/jdcloud/lima-probe-push.timer
[Unit]
Description=Push JDCloud probe results to LiMa every 5 minutes

[Timer]
Unit=lima-probe-push.service
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
